# script/startup_scripts/GENTOO
#!/sbin/runscript
# Copyright 1999-2004 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2

depend() {
	need net
}

start() {
        ebegin "Starting snag"
        /opt/snag/bin/snagw.pl --nosyslog 
        eend $?
}

stop() {
        ebegin "Stopping snag"

        awkfield2='$2'
        procs=(`ps -ef | awk "/(_snags.pl|snagc.pl)/ && ! /awk/ "| awk "{ print $awkfield2; }"`)
                        
        for proc in ${procs[@]}
        do
                ebegin " Killing ${proc}: "
                kill ${proc}
        done
        eend $?
}

restart() {
        svc_stop
        svc_start
}
